# packaging/deb/mayfly_0.0.1-1_arm64/DEBIAN/prerm
systemctl stop mayfly
systemctl disable mayfly
